# extras/virtled.service
[Unit]
Description=Virtual Led Management service
After=network.target

[Service]
ExecStart=/usr/bin/virtled &

[Install]
WantedBy=default.target
